# etc/profile.d/xdg-dektop-portal-kde.sh
# Use the KDE file selection dialog in gnome and other gtk apps.
export GTK_USE_PORTAL=1
